# containers/iot/supervisor/rtsp.conf
[program:rtsp-server]
command=/start-rtsp.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/rtsp-server.err.log
stdout_logfile=/var/log/rtsp-server.out.log