# containers/iot/supervisor/telnet.conf
[program:telnet]
command=/usr/sbin/telnetd -F
autostart=true
autorestart=true
stderr_logfile=/var/log/telnet.err.log
stdout_logfile=/var/log/telnet.out.log
